# etc/config/includes.chroot/etc/skel/.config/systemd/user/appimagelauncherd.service
[Unit]
Description=AppImageLauncher daemon

[Service]
ExecStart=/user/share/appimage-applications/appimagelauncher-lite.AppImage appimagelauncherd
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
